# templates/default/etc/profile.d/docker-compose.sh.erb
export PATH=/opt/docker-compose/bin:$PATH
